# glances/rootfs/etc/services.d/glances/run
#!/usr/bin/env bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: Glances
# Runs Glances
# ==============================================================================
declare -a options=()
declare -a exporters=()

bashio::log.info 'Starting Glances...'

options+=(-C /etc/glances.conf)
options+=(-w)
options+=(--port 61209)
options+=(--bind 127.0.0.1)
options+=(--time "$(bashio::config 'refresh_time')")

if bashio::config.false 'process_info'; then
    options+=(--disable-process)
fi

if bashio::config.true 'influxdb.enabled'; then
    if bashio::config.equals 'influxdb.version' '2'; then
        exporters+=(influxdb2)
    else
        exporters+=(influxdb)
    fi
fi

# Glances v4 does not support --export in webserver mode (-w).
# Run a separate lightweight instance in quiet/standalone mode for exports.
if (( ${#exporters[@]} > 0 )); then
    bashio::log.info 'Starting Glances exporter...'
    declare -a exporter_options=(
        -C /etc/glances.conf
        --quiet
        --time "$(bashio::config 'refresh_time')"
        --export "$(IFS=,; echo "${exporters[*]}")"
    )
    if bashio::config.false 'process_info'; then
        exporter_options+=(--disable-process)
    fi
    glances "${exporter_options[@]}" &
    EXPORTER_PID=$!
    bashio::log.info "Glances exporter started (PID: ${EXPORTER_PID})"
fi

if bashio::debug; then
    options+=(--debug)
fi

# When an exporter is running, both processes run in the background so we can
# detect if either one crashes (wait -n exits on the first child to terminate).
# SIGTERM from HA (addon stop) is forwarded to both children for clean shutdown.
# Without an exporter, the simple exec path is used.
if [[ -n "${EXPORTER_PID:-}" ]]; then
    glances "${options[@]}" &
    WEBSERVER_PID=$!
    bashio::log.info "Glances web server started (PID: ${WEBSERVER_PID})"

    trap 'kill "${WEBSERVER_PID}" "${EXPORTER_PID}" 2>/dev/null' TERM INT
    wait -n
    bashio::log.fatal "A Glances process exited unexpectedly, stopping container..."
    kill "${WEBSERVER_PID}" "${EXPORTER_PID}" 2>/dev/null
    exit 1
else
    exec glances "${options[@]}"
fi
